# zhangjn_hw2/zhangjnhw2/toggle.service
[Unit]
Description= "tog"
Required= flash.service
After= flash.service


[Service]
ExecStart= /toggle.sh
StartLimitInterval= 10s
StartLimiBust= 30